# scripts/daz.service
[Unit]
Description=Daz Chat Bot for Cytube
Documentation=https://github.com/hildolfr/daz
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=daz
Group=daz
WorkingDirectory=/opt/daz

# ExecStart will be configured via override.conf or environment file
ExecStart=/opt/daz/daz

Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/daz/logs

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Environment
Environment="DAZ_LOG_DIR=/opt/daz/logs"
# Load environment variables from file
EnvironmentFile=-/etc/daz/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daz

[Install]
WantedBy=multi-user.target